# souk/deploy/cron.d/install.sh
#!/bin/bash

BASEDIR="$( cd "$( dirname "${BASH_SOURCE[0]}" )" && pwd )"

# cp "${BASEDIR}/cron_task_here" /etc/cron.d/
echo "cron files installed"
